# systemd/octocam-wifi-setup.service
[Unit]
Description=OctoCam first-boot Wi-Fi setup AP
Wants=NetworkManager.service
After=NetworkManager.service
Before=octocam-web.service

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=OCTOCAM_WIFI_CACHE_PATH=/var/lib/octocam/wifi-networks.json
ExecStart=/usr/local/bin/octocam-web --wifi-setup

[Install]
WantedBy=multi-user.target
